# infra/systemd/pawd.service
# PAW Blockchain Node - systemd service file
#
# Installation:
#   sudo cp pawd.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable pawd
#   sudo systemctl start pawd
#
# Configuration:
#   Edit /etc/paw/pawd.env for environment variables
#   Edit ~/.paw/config/app.toml and config.toml for node config

[Unit]
Description=PAW Blockchain Node
Documentation=https://github.com/decristofaroj/paw
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=500
StartLimitBurst=5

[Service]
# User and group (create with: useradd -r -s /sbin/nologin paw)
User=paw
Group=paw

# Environment
EnvironmentFile=-/etc/paw/pawd.env
Environment=HOME=/home/paw
Environment=PAW_HOME=/home/paw/.paw

# Working directory
WorkingDirectory=/home/paw

# Start command
ExecStart=/usr/local/bin/pawd start \
    --home ${PAW_HOME} \
    --log_level ${LOG_LEVEL:-info} \
    --minimum-gas-prices ${MIN_GAS_PRICES:-0.001upaw}

# Restart policy
Restart=on-failure
RestartSec=10
RestartPreventExitStatus=137

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/paw/.paw
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
MemoryDenyWriteExecute=true
RestrictSUIDSGID=true
RestrictNamespaces=true
RestrictRealtime=true
LockPersonality=true

# Capabilities
CapabilityBoundingSet=CAP_NET_BIND_SERVICE
AmbientCapabilities=CAP_NET_BIND_SERVICE

# Resource limits
LimitNOFILE=65535
LimitNPROC=8192
MemoryMax=16G
TasksMax=4096

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=pawd

# Kill settings
KillMode=process
KillSignal=SIGTERM
TimeoutStopSec=120

[Install]
WantedBy=multi-user.target
