# validator-onboarding/quickstart-pack/pawd.service
[Unit]
Description=PAW Validator (paw-mvp-1)
After=network-online.target
Wants=network-online.target

[Service]
User=paw
Group=paw
Environment=DAEMON_HOME=/var/lib/paw
Environment=DAEMON_NAME=pawd
Environment=CHAIN_ID=paw-mvp-1
ExecStart=/usr/local/bin/pawd start --home ${DAEMON_HOME}
ExecStop=/bin/kill -s SIGTERM $MAINPID
Restart=on-failure
RestartSec=5
LimitNOFILE=65535
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
